# packages/integration/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/run/fitb/supervisord.pid
childlogdir=/data/logs

; Unix socket must NOT live under /data: bind-mounting /data from macOS or Windows
; (Docker Desktop) breaks AF_UNIX bind() with EINVAL — supervisord then exits with
; "Cannot open an HTTP server". /run is container-local tmpfs.
[unix_http_server]
file=/run/fitb/supervisor.sock
chmod=0770
chown=foxinthebox:foxinthebox

[supervisorctl]
serverurl=unix:///run/fitb/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── tailscaled (must run as root for NET_ADMIN) ───────────────────────────────
[program:tailscaled]
command=tailscaled --state=/data/data/tailscale/tailscaled.state
user=root
autostart=true
autorestart=true
stdout_logfile=/data/logs/tailscaled.log
stderr_logfile=/data/logs/tailscaled.err
; v0.7.5: cap per-service logs at 10MB × 4 (matches [supervisord]'s own
; logfile cap). Supervisord's per-program defaults are 50MB × 10 — too
; generous for a self-hosted container where /data is on the user's host
; disk and a chatty service could silently fill it.
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=10

; ── tailscale operator watchdog (FITB#127) ────────────────────────────────────
; tailscaled clears its OperatorUser preference on the NeedsLogin transition
; that follows a tailnet key expiry. Once cleared, `tailscale up` from the
; foxinthebox-owned webui returns "Access denied: checkprefs access denied".
; The entrypoint's one-shot `tailscale set --operator=foxinthebox` only
; covers boot — this loop covers the entire container lifetime, including
; key-expiry recovery. Idempotent (no-op when already set). Runs as root
; because OperatorUser can only be (re)set with root or current-operator
; authority — and we just lost the latter when it was cleared.
[program:ts-operator-watchdog]
command=/app/scripts/tailscale-operator-watchdog.sh
user=root
autostart=true
autorestart=true
stdout_logfile=/data/logs/ts-operator-watchdog.log
stderr_logfile=/data/logs/ts-operator-watchdog.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=15

; ── qdrant ───────────────────────────────────────────────────────────────────
[program:qdrant]
command=/app/qdrant/qdrant --config-path /data/config/qdrant.yaml
user=foxinthebox
autostart=true
autorestart=true
stdout_logfile=/data/logs/qdrant.log
stderr_logfile=/data/logs/qdrant.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=20

; ── hermes gateway ────────────────────────────────────────────────────────────
; Wrapped with run-with-env.sh so `supervisorctl restart` picks up keys the
; onboarding wizard wrote to /data/config/hermes.env without a container restart.
[program:hermes-gateway]
command=/app/scripts/run-with-env.sh python -m hermes_cli.main gateway run --replace
user=foxinthebox
autostart=true
autorestart=true
stdout_logfile=/data/logs/hermes-gateway.log
stderr_logfile=/data/logs/hermes-gateway.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/app",PYTHONPATH="/data/apps/hermes-agent",PATH="/usr/local/bin:/usr/bin:/bin",HERMES_HOME="/data/data/hermes",BRAVE_API_KEY="__BRAVE_API_KEY__"
priority=30

; ── hermes webui ──────────────────────────────────────────────────────────────
[program:hermes-webui]
command=/app/scripts/run-with-env.sh python /data/apps/hermes-webui/server.py
user=foxinthebox
autostart=true
autorestart=true
stdout_logfile=/data/logs/hermes-webui.log
stderr_logfile=/data/logs/hermes-webui.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; QA fix v0.4.7-WaveI: HERMES_WEBUI_STATE_DIR was unset, so config.py
; defaulted to HOME/.hermes/webui = /app/.hermes/webui. /app is INSIDE
; the image and is wiped on every container restart — meaning every
; setting persisted via /api/settings (theme, hostname_prompted,
; local_fallback_enabled, the 6 tailscale_* power-user flags) was
; resetting to defaults on every restart. Verified during v0.4.7 QA
; Phase 0 against a real container — set 6 power-user fields via
; /api/settings, restarted, all 6 reset to defaults. Hostname
; survived because that's stored in /data/config/hermes.env (a
; different code path), but every other persisted setting was lost.
; This bug has been latent since the project's beginning. Fix: pin
; STATE_DIR explicitly to /data/state/webui, which is on the
; persistent volume.
; QA fix v0.4.7-WaveJ: HERMES_HOME was set on the gateway
; (`/data/data/hermes`) but NOT on the webui. Webui's use_model() and
; provider-key save paths called _get_config_path() which fell back to
; HOME/.hermes/config.yaml = /app/.hermes/config.yaml — INSIDE the
; image, NOT on the /data volume. Result: every model switch
; (Ollama use-model, OpenRouter key save, etc.) wrote to a file the
; gateway never reads. The gateway reads /data/data/hermes/config.yaml
; per its own HERMES_HOME. Net effect: the user could "switch model"
; via Settings, see the API return ok=true, but the next chat would
; route through whatever stale config was in the gateway's view, AND
; the change wouldn't persist across container restart. Verified
; during v0.4.7 QA: planted stale azure_endpoint into
; /data/config/hermes.yaml, called /api/ollama/use-model,
; /data/data/hermes/config.yaml didn't even exist while
; /app/.hermes/config.yaml had the new model block.
environment=HOME="/app",PYTHONPATH="/data/apps/hermes-webui",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_AGENT_DIR="/data/apps/hermes-agent",HERMES_WEBUI_STATE_DIR="/data/state/webui",HERMES_HOME="/data/data/hermes",ONBOARDING_PATH="/data/config/onboarding.json",PATH="/usr/local/bin:/usr/bin:/bin"
priority=40

; ── llama-server (local AI fallback) ─────────────────────────────────────────
; Bundled llama.cpp llama-server, OpenAI-compat HTTP at 127.0.0.1:8643. Issue
; #9 — provides the local model when the user enables "Local fallback" in
; Settings → Providers AND a GGUF model has been downloaded via the
; /api/local-models manager (#10).
;
; autostart=false: no idle RAM cost for users who don't opt into local
; fallback. api/local_fallback.py issues `supervisorctl start llama-server`
; once a model file is present and the toggle is on. With
; --sleep-idle-seconds 60 the process stays supervised but the model
; weights unload after 60s of no requests — ~50 MB at idle, ~3.5 GB while
; serving. (See #9 Phase 1 research: docs/decisions captured on issue.)
;
; -m points at /data/models/microsoft_Phi-4-mini-instruct-Q4_K_M.gguf
; (matches the canonical filename the download manager writes). When that
; file is missing, the program will fail-fast on start; api/local_fallback
; declines to start it in that case so we don't spam restart loops.
[program:llama-server]
command=/app/llama-cpp/llama-server -m /data/models/microsoft_Phi-4-mini-instruct-Q4_K_M.gguf --host 127.0.0.1 --port 8643 -c 4096 -t 4 --sleep-idle-seconds 60
user=foxinthebox
autostart=false
autorestart=true
startretries=2
stdout_logfile=/data/logs/llama-server.log
stderr_logfile=/data/logs/llama-server.err
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=50
