# packages/integration/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/data/run/supervisord.pid
childlogdir=/data/logs

[unix_http_server]
file=/data/run/supervisor.sock
chmod=0770
chown=foxinthebox:foxinthebox

[supervisorctl]
serverurl=unix:///data/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── tailscaled (must run as root for NET_ADMIN) ───────────────────────────────
[program:tailscaled]
command=tailscaled --state=/data/data/tailscale/tailscaled.state
user=root
autostart=true
autorestart=true
stdout_logfile=/data/logs/tailscaled.log
stderr_logfile=/data/logs/tailscaled.err
priority=10

; ── qdrant ───────────────────────────────────────────────────────────────────
[program:qdrant]
command=/app/qdrant/qdrant --config-path /data/config/qdrant.yaml
user=foxinthebox
autostart=true
autorestart=true
stdout_logfile=/data/logs/qdrant.log
stderr_logfile=/data/logs/qdrant.err
priority=20

; ── hermes gateway ────────────────────────────────────────────────────────────
[program:hermes-gateway]
command=python -m hermes_cli.main gateway run --replace
user=foxinthebox
autostart=true
autorestart=true
stdout_logfile=/data/logs/hermes-gateway.log
stderr_logfile=/data/logs/hermes-gateway.err
environment=HOME="/app",PYTHONPATH="/data/apps/hermes-agent",PATH="/usr/local/bin:/usr/bin:/bin"
priority=30

; ── hermes webui ──────────────────────────────────────────────────────────────
[program:hermes-webui]
command=python /data/apps/hermes-webui/server.py
user=foxinthebox
autostart=true
autorestart=true
stdout_logfile=/data/logs/hermes-webui.log
stderr_logfile=/data/logs/hermes-webui.err
environment=HOME="/app",PYTHONPATH="/data/apps/hermes-webui",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_AGENT_DIR="/data/apps/hermes-agent",ONBOARDING_PATH="/data/config/onboarding.json",PATH="/usr/local/bin:/usr/bin:/bin"
priority=40
