# packages/deb/templates/foxinthebox-updater.service.tmpl
[Unit]
Description=Fox in the Box updater
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'apt-get install --only-upgrade -y foxinthebox && systemctl restart foxinthebox && rm -f __DATA_DIR__/update.trigger'
